# capwords-app.conf
[program:gunicorn-app]
command = /usr/local/bin/gunicorn -b 0.0.0.0:8000 --chdir /mnt/capitolwords/capitolweb capitolweb.wsgi
[program:nginx-app]
command = /usr/sbin/nginx
